# container-files/etc/supervisord.d/php-fpm.conf
[program:php-fpm]
command = php-fpm
autorestart = true
stdout_logfile = /var/log/nginx/php-fpm.log
stderr_logfile = /var/log/nginx/php-fpm.log
